# ushare.crontab
#!/bin/bash
#Asking ushare to update its content
killall -1 ushare
